# unifi/context/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:unifi]
command=/entrypoint.sh /etc/init.d/unifi start
